# root/etc/services.d/ipfs/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc ipfs daemon --enable-pubsub-experiment --enable-namesys-pubsub